# backend/conf/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log

[program:backend-server]
directory=/root/app/
command=python3 server.py
autorestart=true